# overlay/etc/supervisor/conf.d/vnc.conf

[program:x11vnc]
autostart=true
autorestart=true
priority=30
user=%(ENV_USER)s
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -rfbport 5900 -shared -forever
stopsignal=QUIT
stdout_logfile=/home/%(ENV_USER)s/.cache/log/x11vnc.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/home/%(ENV_USER)s/.cache/log/x11vnc.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7

[program:novnc]
autostart=true
autorestart=true
priority=30
user=%(ENV_USER)s
command=/opt/noVNC/utils/launch.sh --vnc localhost:5900 --listen 8083
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=QUIT
stdout_logfile=/home/%(ENV_USER)s/.cache/log/novnc.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/home/%(ENV_USER)s/.cache/log/novnc.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
