# overlay/etc/supervisor/conf.d/xfce4.conf

[program:xfce4]
autostart=true
autorestart=true
priority=50
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
command=/usr/bin/startxfce4
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s"
stopsignal=QUIT
stdout_logfile=/home/%(ENV_USER)s/.cache/log/xfce4.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/home/%(ENV_USER)s/.cache/log/xfce4.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
