# overlay/etc/supervisor/conf.d/xorg.conf

[program:xorg]
autostart=false
autorestart=true
priority=20
user=root
command=/usr/bin/start-xorg.sh
environment=DISPLAY="%(ENV_DISPLAY)s",DISPLAY_DPI="%(ENV_DISPLAY_DPI)s",XDG_RUNTIME_DIR="/tmp/.X11-unix"
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/xorg.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/xorg.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
