# overlay/etc/supervisor.d/dbus.ini

[program:dbus]
priority=10
autostart=true
autorestart=true
directory=/
user=%(ENV_USER)s
command=dbus-daemon --config-file=/usr/share/dbus-1/system.conf --nofork --nopidfile
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/var/log/supervisor/dbus.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/dbus.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
