# overlay/etc/supervisor.d/plugins.ini

[program:plugins]
priority=50
autostart=false
autorestart=false
user=%(ENV_USER)s
directory=/home/%(ENV_USER)s
command=/usr/bin/start-plugins.sh
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
stopsignal=INT
stdout_logfile=/home/%(ENV_USER)s/.cache/log/plugins.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/home/%(ENV_USER)s/.cache/log/plugins.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
