# overlay/etc/supervisor.d/udev.ini

[program:udev]
priority=10
autostart=false
autorestart=true
user=root
directory=/
command=/usr/bin/start-udev.sh
stopsignal=INT
stdout_logfile=/var/log/supervisor/udev.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/log/supervisor/udev.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
